# ops/launchd/com.example.knowledge-hub.continuous-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.example.knowledge-hub.continuous-sync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/ABSOLUTE/PATH/TO/knowledge-hub/scripts/run_continuous_sync.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/ABSOLUTE/PATH/TO/knowledge-hub</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>KHUB_CONTINUOUS_SYNC_ROOT</key>
    <string>/ABSOLUTE/PATH/TO/external-storage/knowledge_os</string>
  </dict>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>10</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/ABSOLUTE/PATH/TO/khub-logs/launchd.continuous_sync.out.log</string>

  <key>StandardErrorPath</key>
  <string>/ABSOLUTE/PATH/TO/khub-logs/launchd.continuous_sync.err.log</string>
</dict>
</plist>
