# ops/launchd/com.won.knowledge-hub.ops-report.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.won.knowledge-hub.ops-report</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/won/Desktop/allinone/knowledge-hub/scripts/run_ops_report.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/won/Desktop/allinone/knowledge-hub</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/won/.khub/logs/launchd.ops_report.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/won/.khub/logs/launchd.ops_report.err.log</string>
</dict>
</plist>
