# ops/launchd/com.example.knowledge-hub.daily-source-quality.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.example.knowledge-hub.daily-source-quality</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>~/.khub/bin/run_daily_source_quality_launchd.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string><repo-root></string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>10</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>~/.khub/logs/launchd.daily_source_quality.out.log</string>

  <key>StandardErrorPath</key>
  <string>~/.khub/logs/launchd.daily_source_quality.err.log</string>
</dict>
</plist>
